# sync/zuanfeng-sync.service
[Unit]
Description=Zuanfeng Deploy Auto-Sync to Jetson Device
Documentation=https://github.com/anomalyco/zuanfeng-deploy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/lsyncd /home/rec/.config/zuanfeng-sync/lsyncd.conf.lua
Restart=on-failure
RestartSec=5

# Environment (can be overridden in ~/.config/zuanfeng-sync/env)
EnvironmentFile=-/home/rec/.config/zuanfeng-sync/env

[Install]
WantedBy=default.target
